# playbook/roles/prometheus_exporter/templates/host_metrics_exporter.service.j2
# /etc/systemd/system/{{ prometheus_exporter_service_name }}
# Managed by Ansible - DO NOT EDIT MANUALLY

[Unit]
Description=Host Metrics Exporter using Prometheus (via pkgx)
After=network.target

[Service]
Type=simple

ExecStart=pkgx +python@{{ prometheus_exporter_python_version }} python {{ prometheus_exporter_dest_path }}
Restart=on-failure
RestartSec=5
WorkingDirectory=/opt
User={{ ansible_user_id }}

Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
